# files/system/etc/profile.d/custom-branding.sh
#!/bin/bash
export CUSTOM_OS_NAME="Zaravi OS"
